# i3-battery-alert.service
[Unit]
Description=I3 Low Level Battery Alert Script

[Service]
ExecStart=i3-battery-alert

[Install]
WantedBy=multi-user.target
